# src/applications/systemd/system/greetd.service.d/vt_palette.conf
[Service]
ExecStartPre=/usr/local/bin/vt_console_palette.sh /dev/tty1
