# src/assets/drop-ins/systemd/user/ollama.service.d/10-environment.conf
[Service]
Environment="OLLAMA_KEEP_ALIVE=-1"
