# src/applications/systemd/deluged.service.d/override.conf
[Unit]
Wants=run-media-israel-EMDRIVE.mount
After=run-media-israel-EMDRIVE.mount
