# stow/systemd/dot-config/systemd/user/ollama.service
[Unit]
Description=Ollama Server

[Service]
ExecStart=/usr/bin/ollama serve
