# stow/systemd/dot-config/systemd/user/ollama.service
[Unit]
Description=Ollama Server

[Service]
ExecStart=/usr/bin/ollama serve
Environment=OLLAMA_HOST=0.0.0.0:11434
